# scripts/com.dustinwicker.denver-rec-scraper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dustinwicker.denver-rec-scraper</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/dustinwicker/projects/denver_rec_centers/scripts/weekly_scrape.sh</string>
    </array>
    
    <!-- Run every Sunday at 6 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/dustinwicker/projects/denver_rec_centers/logs/launchd_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/dustinwicker/projects/denver_rec_centers/logs/launchd_stderr.log</string>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>

